# systemd/dotnet-autodeploy.service
[Unit]
Description=Auto-deploy .NET apps from public_html

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dotnet-autodeploy
