# systemd/dotnet-autodeploy.timer
[Unit]
Description=Run .NET auto-deployer every 30s

[Timer]
OnBootSec=10s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=dotnet-autodeploy.service

[Install]
WantedBy=timers.target
