# systemd/dotnet-autodeploy.service
[Unit]
Description=CyberPanel .NET auto-deploy scanner

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dotnet-autodeploy
