# systemd/dotnet-app@.service
[Unit]
Description=CyberPanel .NET App for %I
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/dotnet-apps/%i.env
WorkingDirectory=%h
ExecStart=/usr/local/bin/cyberpanel-dotnet-wrapper
Restart=always
RestartSec=5
KillSignal=SIGINT
SyslogIdentifier=dotnet-%i

[Install]
WantedBy=multi-user.target
