# systemd/dotnet-app@.service
[Unit]
Description=.NET App %i
After=network.target

[Service]
EnvironmentFile=-/etc/dotnet-apps/%i.env
# WorkingDirectory is set per-app in /etc/systemd/system/dotnet-app@<app>.service.d/override.conf
ExecStart=/usr/bin/dotnet ${WORKING_DIR}/${EXEC_DLL}
Restart=always
RestartSec=5
User=www-data
Group=www-data
SyslogIdentifier=%i

[Install]
WantedBy=multi-user.target
